# crontab
# DDNS定时更新任务
# 从每天0点开始，每4小时的10分更新一次
# 执行时间: 00:10, 04:10, 08:10, 12:10, 16:10, 20:10

10 */4 * * * /app/ddnspod.sh >> /var/log/ddns.log 2>&1
